# root/etc/services.d/redis/run
#!/bin/sh
exec 2>&1

echo "Starting Redis..."
exec s6-setuidgid bsd /usr/local/bin/redis-server /config/redis.conf
